# src/i3Grid.service

[Service]
ExecStart=/usr/local/bin/i3GridController
User=juli


[Install]
WantedBy=default.target

